# blue_insta.service
[Unit]
Description=Blue Instagram Bot Service
After=network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/blue_movie/insta_bot
ExecStart=/root/blue_movie/.venv/bin/python /root/blue_movie/insta_bot/instabot3.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PATH=/root/blue_movie/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/root/blue_movie/insta_bot
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
